# SECURITY/99-usb-sound.rules
ACTION=="add", SUBSYSTEM=="usb", RUN+="/bin/usb-beep.sh"
